# systemd/user/default.target
[Unit]
Description=Default user target
Requires=dbus.service
