# systemd/user/graphical.target
[Unit]
Description=Graphical session
Requires=console.target
Requires=xorg.target
Wants=wm.service
