# systemd/user/graphical.target
[Unit]
Description=Graphical session
BindsTo=xorg.service
Wants=console.target
Wants=wm.service
