# base/systemd/user/console.target
[Unit]
Description=Console target
Requires=base.target
